# var/lib/dpkg/info/unit.conffiles
/etc/default/unit
/etc/init.d/unit
/etc/logrotate.d/unit
